# supervisord.conf
[supervisord]
nodaemon=true

[program:node]
command=node /src/app.js

[program:consul]
command=consul agent -join="172.30.0.58" -data-dir /tmp/consul-data -config-dir /etc/consul.d 
